# build/base-build/supervisord.conf
[supervisord]
nodaemon=true

[program:build]
command=/bin/sh /etc/build_cf_base.sh
autorestart=false
