# build/base-build/supervisord-base.conf
[supervisord]
nodaemon=true

[program:build]
command=/bin/sh /etc/build_rf_base.sh
autorestart=false
