# build/rf-worker/supervisord-worker.conf
[supervisord]
nodaemon=true

[program:setup]
command=python setup.py install --user
autorestart=false
directory=/root/temp/restframe
startsecs=0
priority=1

[program:restframe-dbsync]
command=/root/.local/bin/restframe-dbsync
autorestart=false
startretries=10
priority=20
exitcodes=0

[program:restframe-worker]
command=/root/.local/bin/restframe-worker
autorestart=true
startretries=50
priority=999
