# build/faas-worker/supervisord.conf
[supervisord]
nodaemon=true

[program:faas-worker]
command=/root/.local/bin/faas-worker
autorestart=true
startretries=50
priority=0
