# build/faas-worker/supervisord.conf
[supervisord]
nodaemon=true

[program:setup]
command=python setup.py install --user
autorestart=false
directory=/root/temp/cloudframe
startsecs=0
priority=1

[program:faas-worker]
command=/root/.local/bin/faas-worker
autorestart=true
startretries=50
priority=999
